# crontab
*/10 * * * * python3 ~/buy-mac/check.py "MGNA3KH/A"